# deployments/systemd/nkudo-edge.service
[Unit]
Description=n-kudo edge agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/nkudo-edge run --control-plane ${CONTROL_PLANE_URL}
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/nkudo-edge
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
LimitNOFILE=65536
EnvironmentFile=-/etc/nkudo-edge/nkudo-edge.env

[Install]
WantedBy=multi-user.target
